# files/crontab
* * * * * date > /home/moderne/cron.log

